# deploy/waf-honeypot-safeline.service
[Unit]
Description=WAF Honeypot Collector - SafeLine Syslog Receiver
Documentation=https://github.com/wzw57/waf-honeypot-collector
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/waf-honeypot-collector
ExecStart=/home/ubuntu/waf-honeypot-collector/venv/bin/python main.py recv-safeline --config /home/ubuntu/waf-honeypot-collector/config.yaml
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/waf-honeypot-collector.env

# 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/home/ubuntu/waf-honeypot-collector/data /home/ubuntu/waf-honeypot-collector/logs

[Install]
WantedBy=multi-user.target
